# themes/cloak/root/etc/uci-defaults/luci-theme-cloak
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Cloak=/luci-static/cloak
	commit luci
EOF
exit 0
