# themes/anonymizer/root/etc/uci-defaults/luci-theme-anonymizer
#!/bin/sh
uci batch <<-EOF
	set luci.themes.OpenWrt=/luci-static/anonymizer
        commit luci
EOF
	
